# supervisord.conf
[supervisord]
nodaemon=true

[program:sabnzbd]
command=sabnzbdplus -s 0.0.0.0:8080
